# docker/supervisor/webserver/php.conf
[program:php]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
